# deploy/systemd/mydiary-updater.timer.template
[Unit]
Description=Run MyDiary Auto-Updater periodically

[Timer]
OnBootSec=2min
OnUnitActiveSec={{UPDATE_INTERVAL}}
Unit=mydiary-updater.service
Persistent=true

[Install]
WantedBy=timers.target
