# deploy/systemd/mydiary-web.service.template
[Unit]
Description=MyDiary Web (static frontend via API server)
After=network.target mydiary-api.service
# The Express API server now serves dist/ statically on PORT (default 4000).
# This service is kept for compatibility but is intentionally a no-op.
# The frontend is served by mydiary-api.service.

[Service]
Type=oneshot
ExecStart=/bin/true
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
