# deploy/systemd/mydiary-updater.service.template
[Unit]
Description=MyDiary Auto-Updater (git pull + rebuild + restart if new commits)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=root
Group=root
WorkingDirectory={{APP_DIR}}
Environment=NODE_ENV=production
Environment=HOME={{APP_DIR}}
Environment=PATH=/usr/local/bin:/usr/bin:/bin:{{APP_DIR}}/node_modules/.bin
ExecStart={{APP_DIR}}/scripts/update.sh
